# systemd/trending-topics-worker.service
[Unit]
Description=Trending Topics RabbitMQ Worker
After=network.target rabbitmq-server.service
Requires=rabbitmq-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/get_trend_twitter
Environment="PATH=/var/www/get_trend_twitter/venv/bin"
EnvironmentFile=/var/www/get_trend_twitter/.env

# Run worker script
ExecStart=/var/www/get_trend_twitter/venv/bin/python /var/www/get_trend_twitter/worker.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/trending-topics/worker.log
StandardError=append:/var/log/trending-topics/worker-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/trending-topics

# Resource limits (worker needs more resources for ML processing)
LimitNOFILE=65536
LimitNPROC=8192
MemoryLimit=8G

[Install]
WantedBy=multi-user.target
